# Digital Ocean/user-alert-checker.service
[Unit]
Description=VolSpike User Alert Checker
After=network.target

[Service]
Type=simple
User=trader
WorkingDirectory=/home/trader/volume-spike-bot
EnvironmentFile=/home/trader/.volspike.env
ExecStart=/home/trader/volume-spike-bot/.venv/bin/python3 -u /home/trader/volume-spike-bot/user_alert_checker.py
Restart=always
RestartSec=10
StandardOutput=append:/home/trader/volume-spike-bot/logs/user-alert-checker.log
StandardError=append:/home/trader/volume-spike-bot/logs/user-alert-checker.error.log

[Install]
WantedBy=multi-user.target
